# cron.d/daily.in
# This is a sample cron file for rsnapshot. 
# The values used correspond to the examples in /etc/rsnapshot.conf.
# There you can also set the backup points and many other things.
#
# To activate this cron file you have to uncomment the lines below.
# Feel free to adapt it to your needs.

0 */4   * * *           root    /usr/share/${APPLICATION_PACKAGE_NAME}/scripts/reboot.sh
#30 3    * * *           root    /usr/share/${APPLICATION_PACKAGE_NAME}/scripts/reboot.sh
